# build-infra/builder-cleanup.timer
[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
