# build-infra/builder-cleanup.timer
[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
