# build-infra/builder-cleanup.service
[Service]
ExecStart=/usr/local/bin/builder-cleanup
User=user
Type=oneshot
