# build-infra/upload-release-status.service
[Service]
ExecStart={{builder_dir}}/qubes-src/builder-github/upload-release-status
