# build-infra/upload-release-status.timer
[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
